# units/dietsplash-quit.service.in
[Unit]
Description=Terminate dietsplash Boot Screen
Before=getty@tty1.service display-manager.service
After=dev-tty1.device

[Service]
ExecStart=-@bindir@/dietsplashctl 100 ""
Type=oneshot
